# package/utils/tz_mgr/files/tz_mgr.init
#!/bin/sh /etc/rc.common

START=90
USE_PROCD=1

NETWORK_DEV="br-lan"

start_service() 
{
	procd_open_instance
	procd_set_param user root
	procd_set_param netdev "$NETWORK_DEV"
	procd_set_param command /usr/bin/tz_mgr -i "$NETWORK_DEV"
	procd_set_param respawn\
	${respawn_threshold:-3600} \
	${respawn_timeout:-5} ${respawn_retry:-0}
	procd_close_instance
}


