# package/ramips/applications/dialtool2/files/dialtool2.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007-2012 OpenWrt.org

START=60
STOP=85
start(){
    cfg -e |grep -v 'export value' |  awk '{print $2}' | grep "TZ_DIALTOOL2" > /tmp/dialtool2_config
    /bin/dialtool_new -s /tmp/dialtool2_config  &
}

stop(){
    ps | grep dialtool_new | grep -v grep | cut -d' ' -f 2 | xargs kill -9
    rmmod bm_gobi_net_906.ko
}
